# crontab
0 * * * * /root/vdirsyncer.sh  -v >> /daten/vdirsyncer.log 

